# primesign-publish.service
[Unit]
Description=Primesign Publish API - persists config.json for all visitors
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/runlo-workspace/primesign2
ExecStart=/usr/bin/python3 /home/ubuntu/runlo-workspace/primesign2/publish_api.py 9998
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
